# zip/profile.d/zip.csh
#!/bin/csh
setenv UNZIP "-O gb18030"
setenv ZIPINFO "-O gb18030"