# zip/profile.d/zip.sh
#!/bin/sh
export UNZIP="-O gb18030"
export ZIPINFO="-O gb18030"
